# awbot-irk.service
[Unit]
Description=awbot running irkerd
Requires=network.target

[Service]
ExecStart=/usr/bin/awbot-irk

[Install]
WantedBy=multi-user.target
